# root/etc/services.d/jiacrontab_admin/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc nohup /jiacrontab/jiacrontab_admin/jiacrontab_admin &> /jiacrontab/log/jiacrontab_admin.log &