# root/etc/services.d/jiacrontab_admin/run
#!/usr/bin/with-contenv bash

cd /jiacrontab/jiacrontab_admin
s6-setuidgid abc nohup ./jiacrontab_admin &> ../log/jiacrontab_admin.log &