# root/etc/services.d/jiacrontabd/run
#!/usr/bin/with-contenv bash

cd /jiacrontab/jiacrontabd
s6-setuidgid abc nohup ./jiacrontabd &> ../log/jiacrontabd.log &