# root/etc/services.d/jiacrontabd/run
#!/usr/bin/with-contenv bash

nohup /jiacrontab/jiacrontabd/jiacrontabd &> /jiacrontab/log/jiacrontabd.log &